# root/etc/s6-overlay/s6-rc.d/init-mod-ai-tools-packages/run
#!/usr/bin/with-contenv bash

echo "**** Registering AI tools system packages (screen, tmux, fd-find, ripgrep, fzf) ****"
echo "screen" >> /mod-repo-packages-to-install.list
echo "tmux"   >> /mod-repo-packages-to-install.list
echo "fd-find" >> /mod-repo-packages-to-install.list
echo "ripgrep" >> /mod-repo-packages-to-install.list
echo "fzf"    >> /mod-repo-packages-to-install.list
