# root/etc/s6-overlay/s6-rc.d/init-mod-ai-tools-packages/up
/etc/s6-overlay/s6-rc.d/init-mod-ai-tools-packages/run
